# Contracts/testnet/geth.service
[Unit]
Description=Ethereum go client

[Service]
Type=simple
ExecStart=/usr/bin/geth

[Install]
WantedBy=default.target
