# Contracts/testnet/geth.service
[Unit]
Description=Ethereum go client

[Service]
Type=simple
ExecStart=/opt/tc/gethtc

[Install]
WantedBy=default.target
